# backend/scripts/display-chromium.service
# systemd USER unit — starts Chromium inside the graphical session.
#
# Install with: bash scripts/install-service.sh
#
# Manual install:
#   mkdir -p ~/.config/systemd/user
#   cp scripts/display-chromium.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now display-chromium.service
#   sudo loginctl enable-linger $USER   # start at boot before login

[Unit]
Description=Fullscreen Chromium (display frontend)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple

Environment=XAUTHORITY=%h/.Xauthority
Environment=DISPLAY_OZONE=auto
Environment=CHROMIUM_LOG=/tmp/display-chromium.log

ExecStart=%h/app-frontend/scripts/launch-display-hdmi.sh

Restart=on-failure
RestartSec=5
SuccessExitStatus=0

[Install]
WantedBy=graphical-session.target
